# discord-bot/com.di.discord-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.di.discord-bot</string>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/node</string>
		<string>/Users/xuweijin/Decision-Intelligence-/discord-bot/bot.js</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/xuweijin/Decision-Intelligence-/discord-bot</string>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>ThrottleInterval</key>
	<integer>10</integer>

	<key>StandardOutPath</key>
	<string>/tmp/di-discord-bot.stdout.log</string>

	<key>StandardErrorPath</key>
	<string>/tmp/di-discord-bot.stderr.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/xuweijin</string>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
</dict>
</plist>
